# pautosurround.service.in
[Unit]
Description=Automatic surround selector for PulseAudio
After=pipewire-pulse.service pulseaudio.service
BindsTo=pipewire-pulse.service pulseaudio.service
ConditionUser=!root

[Service]
ExecStart=@BINDIR@/pautosurround
Restart=on-failure
Type=simple

[Install]
WantedBy=pipewire-pulse.service pulseaudio.service
